# packages/server-core/supervisord.conf
[supervisord]
nodaemon=true

[program:feedless]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/feedless/docker-entrypoint.sh -DFOREGROUND"
